# src/crontab
10 2 * * * source /.env && rcb backup > /proc/1/fd/1

